# scripts/chrome-bridge.service
[Unit]
Description=Chrome Bridge — AI-to-Browser relay
Documentation=https://github.com/YOUR_USERNAME/chrome-bridge
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h/chrome-bridge
ExecStart=/usr/bin/python3 -u %h/chrome-bridge/bridge/server.py
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=append:%h/chrome-bridge/runtime/server.log
StandardError=append:%h/chrome-bridge/runtime/server.log

# Security hardening (optional — uncomment if paranoid)
# NoNewPrivileges=yes
# ProtectSystem=strict
# ProtectHome=read-only

[Install]
WantedBy=default.target
